# services/studio/studio.service
[Unit]
Description=Studio service
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/root/container-svc/container-service-gin
ExecStartPre="/bin/sh /root/container-svc/generate_secret.sh > /root/container-svc/secret.key"
WorkingDirectory=/root/container-svc/

[Install]
WantedBy=multi-user.target